# valentine.service
[Unit]
Description=Valentine v2 Multi-Agent AI Assistant
After=network.target redis.service docker.service
Wants=redis.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/valentine
EnvironmentFile=/opt/valentine/.env
ExecStart=/opt/valentine/venv/bin/python -m valentine.main

# Kill ALL child processes when stopping — prevents zombie polling
KillMode=control-group
KillSignal=SIGTERM
TimeoutStopSec=15

# Wait 5s before restart to let Telegram release the polling session
Restart=always
RestartSec=5

# Resource limits (tuned for Oracle Cloud ARM64)
MemoryMax=4G
CPUQuota=90%

# Security — allow writing to workspace, skills, git, src (for /restart pulls), and /tmp
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/opt/valentine/workspace /opt/valentine/.git /opt/valentine/skills /opt/valentine/src /tmp
PrivateTmp=false

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=valentine

[Install]
WantedBy=multi-user.target
