# valentine.service
[Unit]
Description=Valentine v2 Multi-Agent AI Assistant
After=network.target redis.service docker.service
Wants=redis.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/valentine
EnvironmentFile=/opt/valentine/.env
ExecStart=/opt/valentine/venv/bin/python -m valentine.main

# Kill ALL child processes (agents) when stopping — prevents zombie polling
KillMode=control-group
KillSignal=SIGTERM
TimeoutStopSec=15

Restart=always
RestartSec=10

# Resource limits (tuned for Oracle Cloud 1 OCPU / 6GB ARM64)
MemoryMax=4G
CPUQuota=90%

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/tmp/valentine /opt/valentine/workspace
PrivateTmp=false

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=valentine

[Install]
WantedBy=multi-user.target
